# nanit-healthcheck.service
[Unit]
Description=Nanit RTMP stream health check
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/nanit/nanit-healthcheck.sh
